# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:mariadb]
command=/usr/bin/mysqld_safe
autostart=true
autorestart=true
stderr_logfile=/var/log/mariadb.err.log
stdout_logfile=/var/log/mariadb.out.log

[program:keycloak]
command=/opt/keycloak/bin/kc.sh start-dev --import-realm
autostart=true
autorestart=true
stderr_logfile=/var/log/keycloak.err.log
stdout_logfile=/var/log/keycloak.out.log
environment=KC_DB=dev-file,KC_HTTP_PORT=8080,KC_HOSTNAME_STRICT=false,KEYCLOAK_ADMIN=admin,KEYCLOAK_ADMIN_PASSWORD=admin

[program:nextjs]
command=npm start
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/nextjs.err.log
stdout_logfile=/var/log/nextjs.out.log
environment=NODE_ENV=production,PORT=3010
